# honeybadger/bash/etc/systemd/system/badger-code.service
[Unit]
Description=Badger Hot Java

[Service]
Type=simple
ExecStart=/usr/scripts/java-hot-reload.sh

[Install]
WantedBy=multi-user.target
